# supervisord.d/mysql.ini
[program:mysql]
command=mysqld_safe --skip-networking=off --init-file=/tmp/init.sql
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
startsecs=10
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=20
stdout_logfile=/var/log/mysql.log
